# scheduler/crontab
# QuantOps Scheduler - All times are in KST (Asia/Seoul)
# Container timezone is set to Asia/Seoul

# Financial data crawling - daily at midnight KST
0 0 * * * /usr/local/bin/python /app/financial_crawler.py >> /var/log/cron.log 2>&1

# Stock selection - daily at 06:00 KST
0 6 * * * /usr/local/bin/python /app/stock_selection.py >> /var/log/cron.log 2>&1

# Auto rebalancing - daily at 08:30 KST (market opens at 09:00)
30 8 * * * /usr/local/bin/python /app/autorebalance.py >> /var/log/cron.log 2>&1

# Empty line required at end of crontab
